# scripts/dolt-beads-launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.ojfbot.dolt-beads</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/dolt</string>
        <string>sql-server</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>3307</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/yuri/.beads-dolt</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/dolt-beads.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/dolt-beads.err</string>
</dict>
</plist>
